# auditbot.service
[Unit]
Description=Auditbot Telegram Bot Service
After=network.target

[Service]
Type=simple
User=direct
WorkingDirectory=/home/direct/Auditbot
ExecStart=/home/direct/.local/bin/uv run bot.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
